# jarvis.service
[Unit]
Description=Jarvis voice assistant
After=pipewire.service
Requires=pipewire.service

[Service]
Type=simple
# Unix pipe: clap detector → voice pipeline
ExecStart=/bin/sh -c '%h/JarvisClaw/daemon/target/release/jarvis-listen --threshold 3000 | %h/JarvisClaw/.venv/bin/jarvis-brain'
Environment=OPENAI_API_KEY=
Restart=on-failure
RestartSec=2

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
SupplementaryGroups=audio

[Install]
WantedBy=default.target
